# sys-fs/lvm2/files/lvmetad.initd-2.02.116-r4
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

pidfile="/run/lvmetad.pid"
command="/sbin/lvmetad"
command_args="${LVMETAD_OPTS:=-p ${pidfile}}"
start_stop_daemon_args="--pidfile ${pidfile}"

depend() {
	:
}

# KK: added pre-start check for required directory to avoid error messages
start_pre()
{
	checkpath --directory --owner root:root --mode 0700 /run/lvm || return 1
}
