# net-misc/x11vnc/files/x11vnc-init
#!/sbin/runscript

depend() {
	need net
	after xdm
}

start() {
	ebegin "Starting x11vnc"
	start-stop-daemon -S -q -x /usr/local/sbin/x11vnc-start \
	-m -p /var/run/x11vnc.pid
	eend $?
}

stop() {
	ebegin "Stopping x11vnc"
	start-stop-daemon -K -q -n x11vnc
	eend $?
}
