# net-misc/x11vnc/files/x11vnc-init
#!/sbin/runscript

depend() {
	need net xdm
	after xdm
}

checkconfig() {
	if [ -z "$AUTH" ] || [ -z "$RFBAUTH" ] || [ -z "$RFBPORT" ] || \
	[ -z "$LOGFILE" ] || [ -z "$IP_ADDRESS" ]; then
	eerror "You must set config options in /etc/conf.d/x11vnc first"
	return 1
	fi
}

start() {
	ebegin "Starting x11vnc"
	start-stop-daemon -S -q -n x11vnc -x /usr/bin/x11vnc -- $OPTS
	eend $?
}

stop() {
	ebegin "Stopping x11vnc"
	start-stop-daemon -K -q -n x11vnc
	eend $?
}
