# net-print/cndrvcups-capt-bin/files/ccpd.service
[Unit]
Description=Canon Advanced Printing Technology (CAPT) backend daemon
Requires=cups.service

[Service]
ExecStart=/usr/sbin/ccpd
