# x11-misc/x11vnc/files/x11vnc.service
[Unit]
Description=VNC Server for X11
Requires=graphical.target

[Service]
ExecStart=/usr/bin/x11vnc
