# x11-misc/x11vnc/files/x11vnc-init
#!/sbin/runscript

depend() {
	need net xdm
	after xdm
}

checkconfig() {
	if [ -z "$AUTH" ] || [ -z "$RFBAUTH" ] || [ -z "$RFBPORT" ] || \
	[ -z "$LOGFILE" ] || [ -z "$IP_ADDRESS" ] || [ -z "$DISPLAY" ]; then
	eerror "You must set config options in /etc/conf.d/x11vnc first"
	return 1
	fi
	if [ $AUTH="guess" ]; then
	$OPTIONS="-env FD_XDM=1 $SETTINGS"
	else
	$OPTIONS="$SETTINGS"
	fi
}

start() {
	ebegin "Starting x11vnc"
	start-stop-daemon -S -q -b -m -p /var/run/x11vnc.pid \
	-x /usr/bin/x11vnc -- $OPTIONS -loop -forever
	eend $?
}

stop() {
	ebegin "Stopping x11vnc"
	pid=`cat /var/run/x11vnc.pid`

	# kill parent process
	start-stop-daemon -K -q -p /var/run/x11vnc.pid

	# kill child processes
	if [ -n "$pid" ]; then
	pkill -s $pid
	fi

	eend $?
}
